# stack/config/subversion/system/profile.d/subversion.sh
export SUBVERSION_HOME=/usr/share/subversion

export PATH=${PATH}:${SUBVERSION_HOME}/bin
