# stack/config/jenkins/system/profile.d/ant.sh
export ANT_HOME=/usr/share/ant/apache-ant-1.10.5

export PATH=${PATH}:${ANT_HOME}/bin
