# stack/config/jenkins/system/profile.d/sonarqube.sh
export SONARQUBE_HOME=/usr/share/sonarqube

export PATH=${PATH}:${SONARQUBE_HOME}/bin/linux-x86-64
