# stack/config/nexus/system/profile.d/nexus.sh
export NEXUS_HOME=/usr/share/nexus
export NEXUS_PORT=32181
export NEXUS_VERSION=3.15.1-01
