# stack/config/nexus/system/profile.d/nexus.sh
export NEXUS_VERSION=3.15.1-01
export NEXUS_HOME=/usr/share/nexus/nexus-${NEXUS_VERSION}
export ENV NEXUS_DATA=/usr/share/nexus/data
export NEXUS_PORT=32181
